# deploy/wasabi-backup.timer
# 매일 data.db 백업
#   상태 보기 : systemctl list-timers wasabi-backup
#   지금 실행 : sudo systemctl start wasabi-backup.service

[Unit]
Description=매일 data.db 백업

[Timer]
OnCalendar=daily
# 파이가 꺼져 있어서 놓친 백업은 다음 부팅 때 실행한다
Persistent=true
# 부팅 직후 몰리지 않게 0~10분 무작위 지연
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
